# packages/core/host/networking/dhcpcd/overlays/etc/systemd/system/dhcpcd.service.d/10-restart.conf
[Service]
# This is needed to allow hotplugging of usb0:
Restart=always
