# packages/core/host/planktoscope/gpio-init/overlays/usr/lib/systemd/system/planktoscope-org.init-gpio-steppers.service
[Unit]
Description=Release GPIO stepper motors on startup
Before=planktoscope-org.device-backend.controller-adafruithat.service
Before=planktoscope-org.device-backend.controller-planktoscopehat.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/release-gpio-steppers
Restart=no

[Install]
WantedBy=multi-user.target
